# com.contentforge.daily.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.contentforge.daily</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/z/Documents/work/content-forge-ai/run_and_commit.sh</string>
    </array>

    <!-- 每天上午9点运行 -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- 如果错过执行时间（如系统休眠），在唤醒后立即运行 -->
    <key>RunAtLoad</key>
    <false/>

    <!-- 标准输出和错误输出到日志文件 -->
    <key>StandardOutPath</key>
    <string>/Users/z/Documents/work/content-forge-ai/logs/launchd_stdout.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/z/Documents/work/content-forge-ai/logs/launchd_stderr.log</string>

    <!-- 设置环境变量 -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PYTHONPATH</key>
        <string>/Users/z/Documents/work/content-forge-ai</string>
        <key>CONTENT_FORGE_MODE</key>
        <string>auto</string>
    </dict>

    <!-- 确保在系统空闲时运行 -->
    <key>LowPriorityBackgroundIO</key>
    <true/>
    <key>Nice</key>
    <integer>1</integer>
</dict>
</plist>
